# systemd/user/newsboat-reload.service
[Unit]
Description=Reload feeds for Newsboat

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
